# home/.config/systemd/user/checkup.timer
[Unit]
Description=Run update check periodically and on boot

[Timer]
OnBootSec=15min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
